# containers/apache2-7.1/rootfs/etc/cron.d/default
# Bitrix events
* * * * * php -f /var/www/html/local/php_interface/cron_events.php >> /var/log/cron.log
